# misc/systemd/greengrass-lite.target
[Unit]
Description=AWS IoT Greengrass runtime for constrained devices.
AllowIsolate=yes

[Install]
WantedBy=multi-user.target
