# modules/tesd/unit/ggl.aws_iot_tes.socket.in
[Unit]
PartOf=greengrass-lite.target
Before=greengrass-lite.target
StartLimitInterval=20

[Install]
WantedBy=greengrass-lite.target

[Socket]
ListenStream=/run/greengrass/aws_iot_tes
FileDescriptorName=aws_iot_tes
Service=ggl.core.tesd.service
SocketUser=@GGL_SYSTEMD_SYSTEM_USER@
SocketGroup=@GGL_SYSTEMD_SYSTEM_GROUP@
Backlog=10
SocketMode=0660

[Install]
WantedBy=sockets.target
