# ggpubsubd/unit/ggl.gg_pubsub.socket.in
[Unit]
PartOf=greengrass-lite.target
Before=greengrass-lite.target

[Install]
WantedBy=greengrass-lite.target

[Socket]
ListenStream=/run/greengrass/gg_pubsub
FileDescriptorName=gg_pubsub
Service=ggl.core.ggpubsubd.service
SocketUser=@GGL_SYSTEMD_SYSTEM_USER@
SocketGroup=@GGL_SYSTEMD_SYSTEM_GROUP@
Backlog=10
SocketMode=0660

[Install]
WantedBy=sockets.target
