# modules/ggipcd/unit/ggl.gg-ipc.socket.socket.in
[Unit]
PartOf=greengrass-lite.target
Before=greengrass-lite.target
StartLimitInterval=20

[Install]
WantedBy=greengrass-lite.target

[Socket]
ListenStream=/var/lib/greengrass/gg-ipc.socket
FileDescriptorName=gg-ipc.socket
Service=ggl.core.ggipcd.service
SocketUser=@GGL_SYSTEMD_SYSTEM_USER@
SocketGroup=@GGL_SYSTEMD_SYSTEM_GROUP@
Backlog=10
SocketMode=0666

[Install]
WantedBy=sockets.target
