# ggipcd/unit/ggl.ipc_component.socket.in
[Unit]
PartOf=greengrass-lite.target
Before=greengrass-lite.target

[Install]
WantedBy=greengrass-lite.target

[Socket]
ListenStream=/run/greengrass/ipc_component
FileDescriptorName=ipc_component
Service=ggl.core.ggipcd.service
SocketUser=@GGL_SYSTEMD_SYSTEM_USER@
SocketGroup=@GGL_SYSTEMD_SYSTEM_GROUP@
Backlog=10
SocketMode=0660

[Install]
WantedBy=sockets.target
