# bins/tesd/unit/ggl.core.tesd.service.in
[Unit]
StartLimitInterval=20
StartLimitBurst=10
Description=core-bus IoT credential provider
PartOf=greengrass-lite.target

[Install]
WantedBy=greengrass-lite.target

[Service]
Type=exec
ExecStart=@CMAKE_INSTALL_PREFIX@/@CMAKE_INSTALL_BINDIR@/@name@
Restart=always
RestartSec=1
# Disallow from having overly-permissive capabilities
CapabilityBoundingSet=~CAP_SYS_ADMIN ~CAP_SYS_PTRACE
User=@GGL_SYSTEMD_SYSTEM_USER@
Group=@GGL_SYSTEMD_SYSTEM_GROUP@
WorkingDirectory=/var/lib/greengrass

[Install]
WantedBy=multi-user.target

[Unit]
After=ggl.core.ggconfigd.service
After=ggl.core.iotcored.service
# Wait for network devices to be up
After=network.target
# Wait for NTP time
After=systemd-time-wait-sync.service
